# jfreeman-pritunl/templates/mongodb-backup.timer
[Unit]
Description=Run mongodb backup script nightly

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target